# files/crontab
SHELL=/bin/sh
PATH=/sbin:/bin:/usr/sbin:/usr/bin:/usr/local/sbin:/usr/local/bin

*/5 * * * * php -f /usr/local/www/nextcloud/cron.php
